# roles/backup/templates/backup.timer.j2
[Install]
WantedBy=multi-user.target

[Timer]
OnCalendar={{ backup_calendar }}
Unit=backup.service